# archive/tower/cfg/rc.local
#!/bin/sh -e
#
# rc.local
#
# This script is executed at the end of each multiuser runlevel.
# Make sure that the script will "exit 0" on success or any other
# value on error.
#
# In order to enable or disable this script just change the execution
# bits.
#
# By default this script does nothing.

#modprobe -r pcspkr
#ipmitool -H 192.168.1.10 -U ADMIN -P ADMIN sensor thres FAN4 lower 100 200 300
vnc4server -geometry 1024x768 -depth 24 :1 &

exit 0
